# sqli/crontab
* * * * * root /evrs/sys-auto/sys-auto
# A blank line
